# dot_config/systemd/user/taskback.service
[Unit]
Description=Export Tasks to json

[Service]
ExecStart=/home/wayne/.local/bin/backtask

[Install]
WantedBy=default.target
