# dot_config/systemd/user/taskback.timer.d/override.conf
[Unit]
Description=Export Tasks to Json

[Timer]
OnBootSec=15m
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
